# debian/OLD/cron.d
#
# Regular cron jobs for the freeside package
#
0 0	* * *	freeside	/usr/bin/freeside-daily fs_daily
