# debian/freeside.postinst
#!/bin/sh

#probably not needed with systemd
/usr/sbin/update-rc.d freeside defaults 23 01

chown -R freeside /usr/local/etc/freeside
rm -fr  /usr/local/etc/freeside/masondata/*

exit 0

