# debian/postinst
#!/bin/sh

chown -R freeside /usr/local/etc/freeside
/usr/sbin/update-rc.d freeside defaults 23 01
exit 0

