# debian/postinst
#!/bin/sh

chown -R freeside /usr/local/etc/freeside
/usr/sbin/update-rc.d freeside defaults 23 01
rm -fr  /usr/local/etc/freeside/masondata/*

exit 0

